# deploy/com.palmeiras.refresh.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.palmeiras.refresh</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/python3</string>
        <string>/Users/rodrigomelo/Documents/projects/palmeiras-agenda/apps/web/scripts/refresh-all.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/rodrigomelo/Documents/projects/palmeiras-agenda/apps/web</string>

    <key>RunAtLoad</key>
    <true/>

    <key>StartInterval</key>
    <integer>900</integer>

    <key>StandardOutPath</key>
    <string>/Users/rodrigomelo/Documents/projects/palmeiras-agenda/.logs/refresh.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/rodrigomelo/Documents/projects/palmeiras-agenda/.logs/refresh.err</string>
</dict>
</plist>
