# deploy/palmeiras-collector.service
[Unit]
Description=Palmeiras Agenda data collector
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=www-data
Group=www-data
WorkingDirectory=/var/www/palmeiras-web
Environment=PYTHONUNBUFFERED=1
Environment=PYTHONDONTWRITEBYTECODE=1
Environment=PALMEIRAS_COLLECTOR_LOCK=/run/palmeiras-collector/collector.lock
EnvironmentFile=/etc/palmeiras-web.env
ExecStart=/var/www/palmeiras-web/.venv/bin/python -m services.collector.palmeiras_collector.runner
RuntimeDirectory=palmeiras-collector
RuntimeDirectoryMode=0755
TimeoutStartSec=20min
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=full
ReadWritePaths=/var/www/palmeiras-web/apps/web/static/crests /run/palmeiras-collector

[Install]
WantedBy=multi-user.target
